# debian/streamsquirrel.service
[Unit]
Description=Stream Squirrel (NDI to AES67 Bridge)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=streamsquirrel
Group=streamsquirrel
WorkingDirectory=/opt/streamsquirrel
EnvironmentFile=-/etc/default/streamsquirrel
ExecStart=/opt/streamsquirrel/.venv/bin/uvicorn app:app --host ${HOST:-0.0.0.0} --port ${PORT:-8080} ${UVICORN_ARGS:-}
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
