# debian/streamsquirrel.prerm
#!/usr/bin/env bash
set -e
systemctl stop streamsquirrel.service || true
systemctl disable streamsquirrel.service || true
